# extra-files/usr/lib/systemd/system/rmtfs.service
[Unit]
Description=Qualcomm remotefs service

[Service]
ExecStart=/usr/bin/rmtfs -r -s -o /lib/firmware/rmtfs
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
